# okgo.conf
description "OkGo Service"
author "Ryan McCaw"

start on (local-filesystems and net-device-up)
stop on runlevel [!2345]

env GIN_MODE=release
env GOMAXPROCS=4
respawn
exec /opt/okgo/okgo
